# fs/var/spool/cron/crontabs/www-data
# Add SuiteCRM's default cron
* * * * * cd $SUITECRM_INSTALL_DIR; php -f cron.php > /dev/null 2>&1